# etc/systemd/system/zfs-load-key@.service.d/wait-for-remote-unlock.conf
[Unit]
# Attempt remote key fetch before falling back to interactive password prompt.
# zfs-load-key@.service already checks keystatus and exits 0 if key is loaded,
# so if zfs-remote-unlock succeeds this instance silently skips the prompt.
# zfs-remote-unlock waits for network internally rather than via systemd ordering
# to avoid a dependency cycle through networking.service -> local-fs.target.
After=zfs-remote-unlock.service
Wants=zfs-remote-unlock.service
